# teleamiguis/cron/cron.cron
#Mins  Hours  Days   Months  Day of the week
*      1     *      *       *       /home/pablo/Desktop/semestre\ 6/desarrollo\ de\ software/Proyecto\ /Teleamiguiss/teleamiguis/cron/task.sh 
